# deb/etc/xdg/autostart/wgupdown-gui.desktop
[Desktop Entry]
Type=Application
Name=WgUpDown GUI
Exec=/usr/local/bin/wgupdown-gui
X-GNOME-Autostart-enabled=true
